# etc/neo4j.cron
*/5 * * * * /usr/bin/neo4j-load-csv.sh /etc/theia/handler_neo4j_1.cfg 1 2>> /etc/theia/neo4j-err-1.txt 1>> /dev/null
